# scripts/pi-gen-stage_zerostick/04-networking/files/90-wifi-devices.rules
ACTION=="add", SUBSYSTEM=="ieee80211", KERNEL=="phy0", RUN+="/usr/sbin/iw phy %k interface add wlan0 type station", RUN+="/usr/sbin/iw phy %k interface add uap0 type __ap"